# extra/rmake.init.in
#!/bin/bash
#
# chkconfig: 345 81 19
# description: Starts the rMake Daemon
#
# Copyright (c) rPath, Inc.
#
# This program is free software: you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation, either version 3 of the License, or
# (at your option) any later version.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with this program.  If not, see <http://www.gnu.org/licenses/>.
#


# Source function library.
. /etc/init.d/functions

RMAKE_RESET="@sbindir@/rmake-server reset"

if [ -e @sysconfdir@/rmake ] ; then
    . @sysconfdir@/rmake
fi

config() {
    cd /etc/rmake;
    @sbindir@/rmake-server config
    return $?
}


start() {
    echo -n $"Starting rMake Server and Repository:"
    cd /etc/rmake;
    ulimit -n 65536
    @sbindir@/rmake-server start
    RETVAL=$?
    [ $RETVAL = 0 ] && touch /var/lock/subsys/rmake && success || failure
    echo
    if [ $RETVAL != 0 ]; then
        echo "rMake: Please check /var/log/rmake/rmake-server.log for details"
    fi
    return $RETVAL
}

startdebug() {
    echo -n $"Starting rMake Server and Repository:"
    cd /etc/rmake;
    ulimit -n 65536
    @sbindir@/rmake-server start -n -d
    RETVAL=$?
    [ $RETVAL = 0 ] && success || failure
    echo
    return $RETVAL
}



stop() {
    echo -n $"Shutting down rMake Server and Repository: "
    cd /etc/rmake;
    @sbindir@/rmake-server stop
    RETVAL=$?
    [ $RETVAL = 0 ] && rm -f /var/lock/subsys/rmake && success || failure
    echo
    return $RETVAL
}

restart() {
    stop
    start
}

reset() {
    status rmake-server > /dev/null
    RUNNING=$?
    if [ $RUNNING -eq 0 ] ; then
        stop
    fi
    resetinternal || return $?
    if [ $RUNNING -eq 0 ] ; then
        start || return $?
    fi
}

resetinternal() {
    echo $"Flushing rMake internal repository and database"
    cd /etc/rmake;
    $RMAKE_RESET
    RETVAL=$?
    echo
    return $RETVAL
}

RETVAL=0

# See how we were called.
case "$1" in
    config)
        config
        ;;
    start)
        start
        ;;
    stop)
        stop
        ;;
    debug)
        startdebug
        ;;
    status)
        status rmake-server
        ;;
    restart|reload)
        restart
        ;;
    reset)
        reset
        ;;
    *)
        echo $"Usage: $0 {start|stop|debug|status|restart|reload|config|reset}"
        exit 1
esac

exit $?
