# extra/rmake.init.in
#/bin/bash
#
# chkconfig: 345 88 12
# description: Starts the rMake Daemon
#
# Copyright (c) rPath, Inc.
#
# This program is free software: you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation, either version 3 of the License, or
# (at your option) any later version.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with this program.  If not, see <http://www.gnu.org/licenses/>.
#


# Source function library.
. /etc/init.d/functions

name=`basename $0`
if [ ${name:0:1} = "S" -o ${name:0:1} = "K" ]
then
    name=${name:3}
fi

prog=rmake-dispatcher
desc="rMake Server and Repository"
pidfile=@lockdir@/${prog}.pid
lockfile=/var/lock/subsys/${name}

PRE_OPTS="--config-file @confdir@/serverrc"
RMAKE_RESET="@sbindir@/$prog reset"

sysconf=@sysconfdir@/$(echo ${name} | sed -e 's/-node//')
[ -e $sysconf ] && . $sysconf

config() {
    cd @confdir@;
    @sbindir@/${prog} config ${PRE_OPTS}
    return $?
}


start() {
    echo -n $"Starting ${desc}: "
    cd @confdir@;
    ulimit -n 65536
    @sbindir@/${prog} start ${PRE_OPTS} ${RMAKE_OPTS}
    RETVAL=$?
    [ $RETVAL = 0 ] && touch $lockfile && success || failure
    echo
    if [ $RETVAL != 0 ]; then
        echo "rMake: Please check @logdir@/${prog}.log for details"
    fi
    return $RETVAL
}


stop() {
    echo -n $"Shutting down ${desc}: "
    cd @confdir@;
    @sbindir@/${prog} stop ${PRE_OPTS}
    RETVAL=$?
    [ $RETVAL = 0 ] && rm -f $lockfile && success || failure
    echo
    return $RETVAL
}

restart() {
    stop
    start
}

reset() {
    status -p $pidfile -l $lockfile ${prog} > /dev/null
    RUNNING=$?
    if [ $RUNNING -eq 0 ] ; then
        stop
    fi
    resetinternal || return $?
    if [ $RUNNING -eq 0 ] ; then
        start || return $?
    fi
}

resetinternal() {
    echo $"Flushing ${desc}: "
    cd @confdir@
    $RMAKE_RESET
    RETVAL=$?
    echo
    return $RETVAL
}

RETVAL=0

# See how we were called.
case "$1" in
    config)
        config
        ;;
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status -p $pidfile -l $lockfile ${prog}
        ;;
    restart|reload)
        restart
        ;;
    reset)
        reset
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|reload|config|reset}"
        exit 1
esac

exit $?
